# opentalon.service.template
[Unit]
Description=OpenCat Web Service
Documentation=https://github.com/ziwei-control/opencat
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%INSTALL_DIR%
ExecStart=%PYTHON% %INSTALL_DIR%/web_server.py --port 6767
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencat

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
